# bin/99-power.rules
SUBSYSTEM=="power_supply", ATTR{online}=="0", RUN+="tlp start"
SUBSYSTEM=="power_supply", ATTR{online}=="1", RUN+="tlp start"
